# deploy/minervini-dashboard.service
[Unit]
Description=Minervini AI — Streamlit Dashboard
Documentation=https://github.com/youruser/minervini_ai
After=network.target minervini-api.service
Wants=network.target
BindsTo=

[Service]
Type=simple

User=@@DEPLOY_USER@@
Group=@@DEPLOY_USER@@

WorkingDirectory=@@PROJECT_DIR@@
EnvironmentFile=@@PROJECT_DIR@@/.env

ExecStart=@@PROJECT_DIR@@/.venv/bin/streamlit run \
    @@PROJECT_DIR@@/dashboard/app.py \
    --server.port 8501 \
    --server.headless true \
    --server.enableCORS false \
    --server.enableXsrfProtection false

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=minervini-dashboard

TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
