# deploy/minervini-daily.service
[Unit]
Description=Minervini AI — Daily SEPA Screen
Documentation=https://github.com/youruser/minervini_ai
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

User=@@DEPLOY_USER@@
Group=@@DEPLOY_USER@@

WorkingDirectory=@@PROJECT_DIR@@
EnvironmentFile=@@PROJECT_DIR@@/.env

ExecStart=@@PROJECT_DIR@@/.venv/bin/python \
    @@PROJECT_DIR@@/scripts/run_daily.py \
    --date today

StandardOutput=journal
StandardError=journal
SyslogIdentifier=minervini-daily

TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
