# deploy/minervini-api.service
[Unit]
Description=Minervini AI — FastAPI Backend
Documentation=https://github.com/youruser/minervini_ai
After=network.target
Wants=network.target

[Service]
Type=simple

User=@@DEPLOY_USER@@
Group=@@DEPLOY_USER@@

WorkingDirectory=@@PROJECT_DIR@@
EnvironmentFile=@@PROJECT_DIR@@/.env

ExecStart=@@PROJECT_DIR@@/.venv/bin/uvicorn \
    api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=minervini-api

TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
